# package/rvi-probe/files/etc/init.d/cloudflared
#!/bin/sh /etc/rc.common
# Simple token-based cloudflared runner for OpenWrt
START=99
USE_PROCD=1

BIN="/usr/bin/cloudflared"
TOKEN_FILE="/etc/cloudflared/token"

start_service() {
    [ -x "$BIN" ] || return 0
    [ -s "$TOKEN_FILE" ] || return 0
    procd_open_instance
    procd_set_param command "$BIN" --no-autoupdate tunnel run --token "$(cat "$TOKEN_FILE")"
    procd_set_param respawn 300 5 5
    procd_close_instance
}

stop_service() { :; }
reload_service() { stop; start; }

