# package/rvi-probe/files/etc/init.d/cloudflared
#!/bin/sh /etc/rc.common
START=90
USE_PROCD=1
PROG=/usr/bin/cloudflared
TOKEN_FILE=/etc/rvi/cloudflared.token

start_service() {
  [ -s "$TOKEN_FILE" ] || { echo "[cloudflared] token missing at $TOKEN_FILE"; return 0; }
  procd_open_instance
  procd_set_param command "$PROG" --no-autoupdate tunnel run --token "$(cat "$TOKEN_FILE")"
  procd_set_param respawn 5 10 5
  procd_close_instance
}
stop_service() { :; }
